# dev-util/dcd/files/dcd-server.service
[Unit]
Description=Dlang Completion Daemon

[Service]
ExecStart=/usr/bin/dcd-server
ExecStop=/usr/bin/dcd-client --shutdown
